# hashirama-chat.service
[Unit]
Description=Hashirama Chat Susanoo UI
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/workspace/hashirama-chat-susanoo
ExecStartPre=/usr/bin/npm --prefix /workspace/hashirama-chat-susanoo run build
ExecStart=/usr/bin/node /workspace/hashirama-chat-susanoo/dist/server.js
Restart=always
RestartSec=3
User=root
Environment=PORT=8090

[Install]
WantedBy=multi-user.target
